# provisioning/ansible/roles/gold-master/templates/ironstone-hostname-fallback.service.j2
[Unit]
Description=Ironstone Hostname Fallback (set hostname from MAC if cloud-init failed)
After=network.target cloud-init.service
Wants=network.target
ConditionPathExists=!/var/lib/cloud/instance/boot-finished

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ironstone-hostname-fallback.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
