# provisioning/ansible/roles/gold-master/templates/k3s-init.service.j2
# =============================================================================
# K3s First-Boot Initialization Service
# =============================================================================
# This service runs once after cloud-init completes to start k3s.
# It expects cloud-init to have written /etc/rancher/k3s/config.yaml
# with the server URL and token from Matchbox.
# =============================================================================

[Unit]
Description=K3s First-Boot Initialization
Documentation=https://k3s.io
After=cloud-init.target network-online.target
Wants=network-online.target
ConditionPathExists=/etc/rancher/k3s/config.yaml
ConditionPathExists=!/var/lib/rancher/k3s/agent/kubelet.kubeconfig

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sleep 10
ExecStart=/usr/local/bin/k3s-init.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
