# provisioning/armbian-build/userpatches/overlay/etc/systemd/system/systemd-networkd-wait-online.service.d/override.conf
[Service]
# Only wait for one interface to be online, not all
# This prevents boot delays when multiple NICs exist but only one is connected
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --any
